# deb/install/postinst
#!/bin/bash

case "${1}" in
  configure)
    echo "Enabling Waggle Agent Power Service"
    systemctl enable waggle-agent-power.service
    ;;
esac
