# deb/install/postinst
#!/bin/bash

case "${1}" in
  configure)
    echo "Enabling Waggle Power-up Service"
    systemctl enable waggle-power-up.service
    ;;
esac
