# deb/install/prerm
#!/bin/bash

systemctl disable waggle-agent-power.service || true
