# deb/install/prerm
#!/bin/bash

systemctl disable waggle-power-up.service || true
